# debian/cryptdisks-enable.upstart
# cryptdisks - enable encrypted block devices
#
# This is a bit of a kludge, we should do this per-block device in a udev
# rule.  But this should work for now.

description	"enable encrypted block devices"

start on stopped udevtrigger

task

script
    . /lib/cryptsetup/cryptdisks.functions

    case "$CRYPTDISKS_ENABLE" in
    [Nn]*)
        exit 1
        ;;
    esac

    INITSTATE="sad"
    do_start
end script
